# php/nginx/fpm.conf
[program:php5-fpm]
command=/opt/fpm.sh
autostart=true
autorestart=true
