# nginx/nginx.conf
[program:nginx]
command=/opt/nginx.sh
numprocs=1
autostart=true
autorestart=true
