# serf/serf.conf
[program:serf]
command=/opt/serf.sh
numprocs=1
autostart=true
autorestart=true
