# mysql/mysql.conf
[program:mysqld]
command=/opt/mysql.sh
numprocs=1
autostart=true
autorestart=true
